# systemd/amdproxy-reverse-tunnel.service
[Unit]
Description=Expose local AMDproxy to a control node via reverse SSH tunnel
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
Environment=AUTOSSH_GATETIME=0
Environment=LOCAL_AMDPROXY_HEALTH=http://127.0.0.1:8082/health
ExecStartPre=/usr/bin/curl --silent --show-error --fail --max-time 5 ${LOCAL_AMDPROXY_HEALTH}
ExecStart=/usr/bin/autossh -M 0 -N -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -o ExitOnForwardFailure=yes -o StrictHostKeyChecking=accept-new -o ControlMaster=no -o ControlPath=none -R 127.0.0.1:8882:127.0.0.1:8082 -R 127.0.0.1:8883:127.0.0.1:8082 CONTROL_NODE_SSH_TARGET
Restart=always
RestartSec=5
KillMode=control-group

[Install]
WantedBy=multi-user.target
